# examples/systemd/red@.service
[Unit]
Description=Red-DiscordBot instance %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/opt/red/YOUR_INSTANCE
EnvironmentFile=-/etc/red/%i.env
ExecStart=/opt/red/YOUR_INSTANCE/.venv-red/bin/redbot %i
Restart=on-failure
RestartSec=15
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
